# config/launchd/ai.openclaw.company-kernel.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.company-kernel.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>__COMPANY_KERNEL_ROOT__/bin/company-daemon</string>
    <string>--once</string>
    <string>--summary</string>
  </array>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>WorkingDirectory</key>
  <string>__COMPANY_KERNEL_ROOT__</string>
  <key>StandardOutPath</key>
  <string>__COMPANY_KERNEL_ROOT__/logs/company-daemon.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>__COMPANY_KERNEL_ROOT__/logs/company-daemon.launchd.err.log</string>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
